# examples/square/supervisord.conf
[supervisord]
logfile:/dev/null
logfile_maxbytes:0

[program:chatserver]
command=./mqtt-fetch-square.py
numprocs=20
process_name=square-backend-%(process_num)s
