# examples/chat/supervisord.conf
[supervisord]
logfile:/dev/null
logfile_maxbytes:0

[program:chatserver]
command=./mqtt-chat.py
numprocs=8
process_name=chat-backend-%(process_num)s
